# supervisord.conf
[program:node]
command = node server.js --port=8080
directory = /home/dotcloud
environment = NODE_ENV=production

[supervisorctl]
serverurl = unix:///var/dotcloud/supervisor.sock
